# kociol7.service
[Unit]
Description=Kociol7

[Service]
ExecStart=/home/pi/kociol7-server/main.py
Type=simple

[Install]
WantedBy=multi-user.target
